# vpnconf.cron
@reboot /etc/openvpn/vpnconf server --server &
@reboot /etc/openvpn/vpnconf client --client
